# roles/hostsctl/files/hostsctl.service
[Unit]
Description=Update /etc/hosts using hostsctl
After=network.target

[Service]
ExecStart=/usr/bin/hostsctl update
