# roles/timezone/files/update-timezone.service
[Unit]
Description=Update /etc/localtime
After=network.target

[Service]
ExecStart=/usr/bin/update-timezone
