# roles/hostsctl/templates/hostsctl.timer.j2
[Unit]
Description=Update /etc/hosts using hostsctl on a schedule

[Timer]
OnCalendar={{ sudoforge_hostsctl_timer_schedule }}
RandomizeDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
