# roles/unbound/templates/roothints.timer.j2
[Unit]
Description=Update root.hints on a schedule

[Timer]
OnCalendar={{ sudoforge_unbound_roothints_timer }}
Persistent=true

[Install]
WantedBy=timers.target
